# service/b.timer
[Unit]
Description=every 24 H

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true

[Install]
WantedBy=timers.target
